# systemd/gh-sync.timer
[Unit]
Description=Sync GitHub copilot-ready issues every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=gh-sync.service

[Install]
WantedBy=timers.target
